# systemd/nicocast.service
[Unit]
Description=NicoCast – Miracast sink for Raspberry Pi Zero 2W
Documentation=https://github.com/nicolasasauer/nicocast-v1
After=network.target wpa_supplicant@wlan0.service

[Service]
Type=simple
# Run as the nicocast user created by the install script.
# Change to 'root' temporarily if wpa_cli requires elevated privileges.
User=nicocast
Group=nicocast

WorkingDirectory=/opt/nicocast

# Pre-start: check operation_mode and stop/start NetworkManager accordingly.
# The '+' prefix runs this command as root regardless of the User= setting.
ExecStartPre=+/opt/nicocast/scripts/nicocast-prestart.sh

ExecStart=/opt/nicocast/venv/bin/python -m nicocast \
    --config /etc/nicocast/nicocast.conf

# Restart on failure with a 5-second delay
Restart=on-failure
RestartSec=5

# Allow the service to use network and GPIO
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nicocast

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=GST_DEBUG=1

[Install]
WantedBy=multi-user.target
